# roles/sslcert/files/certbot-renew.timer
[Unit]
Description=Certbot eene timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
