# roles/sslcert/files/certbot-renew.service
[Unit]
Description=Certbot renew

[Service]
ExecStart=/usr/bin/certbot renew
